# timer/sensor.service
[Unit]
Description = sensor service

[Service]
Type = simple
WorkingDirectory = [path]/Sensor/
ExecStart = python3 [path]/Sensor/sequence.py
